# tools/com.client.erd.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.client.erd.weekly</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__REPO_DIR__/tools/weekly_refresh.sh</string>
  </array>

  <!-- Lunes 07:00 (local) -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>2</integer>
    <key>Hour</key><integer>7</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>__REPO_DIR__/tools/.logs/erd_weekly.out</string>
  <key>StandardErrorPath</key>
  <string>__REPO_DIR__/tools/.logs/erd_weekly.err</string>

  <key>RunAtLoad</key><false/>
</dict>
</plist>
